# systemd/webtoons.socket
[Unit]
Description=webtoons socket

[Socket]
ListenStream=/run/webtoons.sock

[Install]
WantedBy=sockets.target
